# debian/powervr.init
#!/sbin/openrc-run
# Copyright 2020-2021 Maemo Leste
# Distributed under the terms of the GNU General Public License v2

depend() {
	before logger
	after eudev
}

start() {
	ebegin "Starting PowerVR"
	/usr/bin/pvrsrvinit
}

stop() {
	ebegin "Stopping PowerVR (no-op)"
}
