# debian/sgx-ddk-um-tools.postinst
#!/bin/sh

rc-update add powervr sysinit

#DEBHELPER#
